# etc/service/test.timer
[Unit]
Description=test

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target